# configs/home/hauke/.config/systemd/user/rofi-runner.timer
[Unit]
Description=Timer for the rofi runner job.

[Timer]
OnUnitActiveSec=5m
OnBootSec=5m

[Install]
WantedBy=timers.target
